# updater/com.democratdollar.updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.democratdollar.updater</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/claude/projects/apps/democrat-dollar/updater/run-update.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/claude/projects/apps/democrat-dollar/updater</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/claude/.nvm/versions/node/v25.6.1/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>2</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/claude/projects/apps/democrat-dollar/updater/logs/launchd-updater.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/claude/projects/apps/democrat-dollar/updater/logs/launchd-updater-error.log</string>
</dict>
</plist>
